# cron.d/gravityOptimise
00 6 * * * root PATH="$PATH:/usr/local/bin/" gravityOptimise.sh
